# stm32_client.conf
[program:stm32_client]
user=root
command=/usr/bin/stm32_client_bin -d
autostart = true
autorestart = true
